# assistant/glassbox-head-assistant.service
[Unit]
Description=Glassbox Middle Layer (OpenCode)
After=network.target

[Service]
Type=simple
WorkingDirectory=/opt/glassbox-operator
EnvironmentFile=/opt/glassbox-operator/middleware/.env
ExecStart=/usr/bin/python3 -m uvicorn middleware.app:app --host 0.0.0.0 --port 8099
Restart=always
RestartSec=3
User=glassbox

[Install]
WantedBy=multi-user.target
